# debian/ganglios.cron
#
#  update /var/lib/ganglia/xmlcache/ every 2 minutes
#

*/2 * * * * nagios /usr/sbin/ganglia_parser

